# rootfs/etc/cron.d/micha
SHELL=/bin/dash
MAILTO=""
*/5 * * * * root date -u '+\%F \%T' > /etc/fake-hwclock.data
